# cron/cron.daily
#!/bin/bash
REPLACELBPBINDIR/cron_restart.sh
